# x86/backuppc/supervisord.conf
[supervisord]
[program:apache2]
command=/usr/sbin/apache2 -DFOREGROUND

[program:backuppc]
command=gosu backuppc /usr/share/backuppc/bin/BackupPC
