# owncloud/supervisord.conf
[supervisord]
[program:mysqld]
command=mysqld_safe
#numproc=1
#autostart=true
#autorestart=true

[program:nginx]
command=nginx -g "daemon off;"

[program:php-fcgi]
command=php-fpm7.0 -F
